# daemon-emacs/.config/systemd/user/emacs.socket
[Unit]
Description=Emacs Socket

[Socket]
ListenStream=%t/emacs/server
SocketMode=0600
DirectoryMode=0700

[Install]
WantedBy=sockets.target
